# scripts/deployment/systemd/vgindex-backup.timer
[Unit]
Description=Run the VGIndex backup every day at 06:00 UTC

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true
Unit=vgindex-backup.service

[Install]
WantedBy=timers.target
